# deploy/parts-library.service
[Unit]
Description=Parts Library (FastAPI/uvicorn)
After=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/parts-library/backend
ExecStart=/home/pi/parts-library/backend/.venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 8000
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
